# ass1/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

# Make sure java 1.6 is used via 'alternatives --config java'

export M2_HOME=/usr/local/apache-maven/apache-maven-3.2.1/
export M2=$M2_HOME/bin
export JAVA_HOME=/usr/lib/jvm/java

PATH=$JAVA_HOME/bin:$M2:$PATH
export $PATH
